# build/root/etc/s6-overlay/s6-rc.d/init-proton-config/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# Initialize pass
gpg --generate-key --batch /app/protonmail/gpgparams
pass init pass-key

# permissions
lsiown -R abc:abc \
    /config
    

    
