# build/root/etc/s6-overlay/s6-rc.d/init-proton-config/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash
HOME=/config
export HOME

# permissions
lsiown -R abc:abc \
    /config

# Initialize pass
if [ ! -d "/config/.gnupg" ]; then
  s6-setuidgid abc gpg --generate-key --batch /app/protonmail/gpgparams
  s6-setuidgid abc pass init pass-key
fi
    

    
