# build/root/etc/s6-overlay/s6-rc.d/svc-proton-configured/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

if [ ! -d "/config/.password-store/protonmail-credentials/" ]; then
    echo "[BRIDGE] Not configured yet: please run `init-bridge`"
    exit 1
fi


exit 0