# build/root/etc/s6-overlay/s6-rc.d/svc-proton-configured/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash
HOME=/config
export HOME


if s6-setuidgid abc pass; then
    echo "[BRIDGE] Configured: pass has a vault configured."
else
    echo "[BRIDGE] Not configured yet: please run /app/protonmail/init-bridge.sh"
    exit 1
fi
